# docker/cron.conf
[group:cron]
programs=crond
priority=25

[program:crond]
command = cron -f